# mesh/nats/nats-hub.service
[Unit]
Description=NATS hub — constellation fleet (ryzen-work)
After=network.target

[Service]
ExecStart=%h/.local/bin/nats-server -c %h/.config/nats/hub.conf
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
