# ansible/roles/secrets/templates/constellation-secrets-render.service.j2
[Unit]
Description=Render constellation secrets into EnvironmentFiles
Documentation=https://github.com/j0yen/constellation
# Must complete before any secret-consuming units start.
Before=wm-busbridge.service
Before=wm-brain.service
Before=tailscaled.service
Before=nats.service
After=local-fs.target
After=network-pre.target
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart={{ secrets_bin }} audit --role {{ constellation_role | default('') }}
ExecStart=/bin/bash -c ' \
    rm -f {{ secrets_env_dir }}/environment; \
    {% for item in constellation_secrets_map %} \
    val=$(CONSTELLATION_SECRETS_ROOT="{{ secrets_root }}" \
          CONSTELLATION_ROLE="{{ constellation_role | default('') }}" \
          {{ secrets_bin }} get "{{ item.secret }}" 2>/dev/null) && \
        printf "%s=%s\n" "{{ item.env_var }}" "$val" >> {{ secrets_env_dir }}/environment; \
    {% endfor %} \
    chmod 600 {{ secrets_env_dir }}/environment'
# No plaintext in journal
StandardOutput=journal
StandardError=journal
# The rendered file lives in /run (tmpfs); purged on reboot.
RuntimeDirectory=constellation
RuntimeDirectoryMode=0700

[Install]
WantedBy=multi-user.target
WantedBy=wintermute.target
