# cloud/systemd/hub/adopt-cron.service
[Unit]
Description=adopt binary adoption pass
After=network.target

[Service]
ExecCondition=wm-node role hub
Type=oneshot
ExecStart=/bin/bash -c '\
  RUN_ID=$(date +%%Y-%%m-%%d.%%H) ; \
  %h/.local/bin/adopt reconcile 2>&1 | tee /tmp/adopt-cron.log ; \
  %h/.local/bin/adopt apply --execute --with-daemons 2>&1 | tee -a /tmp/adopt-cron.log ; \
  %h/.local/bin/adopt report --run "$RUN_ID" 2>&1 | tee -a /tmp/adopt-cron.log ; \
  echo "[adopt-cron] done $RUN_ID" \
'
StandardOutput=journal
StandardError=journal
