# cloud/systemd/hub/ballast-guard.service
[Unit]
Description=ballast disk SLO guard pass
After=default.target

[Service]
ExecCondition=wm-node role hub
Type=oneshot
ExecStart=%h/.local/bin/ballast-guard-run.sh
StandardOutput=journal
StandardError=journal
