# cloud/systemd/hub/roundtable.service
[Unit]
Description=Roundtable daily session

[Service]
ExecCondition=wm-node role hub
Type=oneshot
ExecStart=%h/.local/bin/roundtable session --with-games
